# roles/ovh-vps-debian-12/files/etc/cron.d/rkunter
0 12 * * * root rkhunter --update
0 13 * * * root rkhunter -c --rwo